# etc-systemd-system/five-minute-timer.target
[Unit]
Description=Five Minute Timer Target
StopWhenUnneeded=yes
