# etc-systemd-system/five-minute-timer.timer
[Unit]
Description=Five Minute Timer

[Timer]
OnCalendar=*:0/15
Unit=five-minute-timer.target

[Install]
WantedBy=timers.target
